# deploy/systemd-docker/stockalert-app.service
[Unit]
Description=Stock Alert Streamlit App (Docker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/app
Environment="PATH=/app/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/app/.env
ExecStart=/app/.venv/bin/streamlit run Home.py --server.port 8501 --server.address 0.0.0.0 --server.headless true
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stockalert-app

# Container-friendly settings (relaxed from bare-metal)
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/app

[Install]
WantedBy=multi-user.target
