# deploy/systemd/stockalert-watchdog.service
[Unit]
Description=Stock Alert Scheduler Watchdog
After=network.target stockalert-scheduler.service stockalert-hourly.service stockalert-futures.service
Wants=stockalert-scheduler.service stockalert-hourly.service stockalert-futures.service

[Service]
Type=simple
User=stockalert
Group=stockalert
WorkingDirectory=/opt/stockalert
Environment="PATH=/opt/stockalert/.venv/bin:/home/stockalert/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="ENVIRONMENT=production"
ExecStart=/home/stockalert/.local/bin/uv run python -m src.services.scheduler_watchdog
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stockalert-watchdog

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/stockalert
PrivateTmp=true

[Install]
WantedBy=multi-user.target
