# deploy/auto-scheduler.service
# systemd unit for the daily/weekly auto scheduler (stock alerts).
# Production: Debian Linux, runs as root. Local: use Docker.
#
# Install on production (Debian):
#   sudo cp deploy/auto-scheduler.service /etc/systemd/system/
#   sudo sed -i 's|/path/to/stock-market-alert-app|/opt/stock-market-alert-app|' /etc/systemd/system/auto-scheduler.service
#   sudo systemctl daemon-reload
#   sudo systemctl start auto-scheduler
#
# Enable on boot:  sudo systemctl enable auto-scheduler
# Logs:           journalctl -u auto-scheduler -f

[Unit]
Description=Daily/weekly auto scheduler (stock alerts, calendar-aware)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/stock-market-alert-app
ExecStart=/path/to/stock-market-alert-app/deploy/start_auto_scheduler.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
