# deploy/futures-scheduler.service
# systemd unit for the futures scheduler (price updates + alert checks).
# Production: Debian Linux, runs as root. Local: use Docker.
#
# Install on production (Debian):
#   sudo cp deploy/futures-scheduler.service /etc/systemd/system/
#   sudo sed -i 's|/path/to/stock-market-alert-app|/opt/stock-market-alert-app|' /etc/systemd/system/futures-scheduler.service
#   sudo systemctl daemon-reload
#   sudo systemctl start futures-scheduler
#
# Enable on boot:  sudo systemctl enable futures-scheduler
# Logs:           journalctl -u futures-scheduler -f

[Unit]
Description=Futures scheduler (price updates and alert checks)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/stock-market-alert-app
ExecStart=/path/to/stock-market-alert-app/deploy/start_futures_scheduler.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
