# deploy/hourly-data-scheduler.service
# systemd unit for the hourly data scheduler (hourly price updates at :05).
# Production: Debian Linux, runs as root. Local: use Docker.
#
# Install on production (Debian):
#   sudo cp deploy/hourly-data-scheduler.service /etc/systemd/system/
#   sudo sed -i 's|/path/to/stock-market-alert-app|/opt/stock-market-alert-app|' /etc/systemd/system/hourly-data-scheduler.service
#   sudo systemctl daemon-reload
#   sudo systemctl start hourly-data-scheduler
#
# Enable on boot:  sudo systemctl enable hourly-data-scheduler
# Logs:           journalctl -u hourly-data-scheduler -f

[Unit]
Description=Hourly data scheduler (price updates at :05 past the hour)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/stock-market-alert-app
ExecStart=/path/to/stock-market-alert-app/deploy/start_hourly_data_scheduler.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
