# deploy/scheduler-watchdog.service
# systemd unit for the scheduler watchdog (monitors daily/weekly and hourly schedulers).
# Production: Debian Linux, runs as root. Local: use Docker.
#
# Install on production (Debian):
#   sudo cp deploy/scheduler-watchdog.service /etc/systemd/system/
#   sudo sed -i 's|/path/to/stock-market-alert-app|/opt/stock-market-alert-app|' /etc/systemd/system/scheduler-watchdog.service
#   sudo systemctl daemon-reload
#   sudo systemctl start scheduler-watchdog
#
# Optional env: SCHEDULER_WATCHDOG_INTERVAL, SCHEDULER_HEARTBEAT_MAX_AGE, WATCHDOG_ENABLE_DISCORD
# Enable on boot:  sudo systemctl enable scheduler-watchdog
# Logs:           journalctl -u scheduler-watchdog -f

[Unit]
Description=Scheduler watchdog (monitors daily/weekly and hourly schedulers)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/stock-market-alert-app
ExecStart=/path/to/stock-market-alert-app/deploy/start_scheduler_watchdog.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
